# dockerfiles/dendrite/dendrite-cdh5/hadoop-yarn-nodemanager.conf
[program:hadoop-yarn-nodemanager]
command=service hadoop-yarn-nodemanager start
environment = JAVA_HOME="/usr/lib/jvm/jdk1.7.0_51"
priority=6
