# dockerfiles/dendrite/hadoop-yarn-resourcemanager.conf
[program:hadoop-yarn-resourcemanager]
command=service hadoop-yarn-resourcemanager start
environment = JAVA_HOME="/usr/lib/jvm/jdk1.7.0_51"
priority=8
