# dockerfiles/dendrite/hadoop-yarn.conf
[program:hadoop-yarn]
command=service hadoop-yarn-proxyserver start
environment = JAVA_HOME="/usr/lib/jvm/jdk1.7.0_51"
