# dockerfiles/dendrite/hbase.conf
[program:hbase]
command=service hbase-master start
environment = JAVA_HOME="/usr/lib/jvm/jdk1.7.0_51"
