# dockerfiles/dendrite/hbase-regionserver.conf
[program:hbase-regionserver]
command=service hbase-regionserver start
environment = JAVA_HOME="/usr/lib/jvm/jdk1.7.0_51"
priority=10
