# dockerfiles/dendrite/npm.conf
[program:npm]
command=npm start
directory=/Dendrite/src/main/nodejs/ungit
