# dockerfiles/dendrite/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd
priority=1
