# dockerfiles/dendrite/rsyslog.conf
[program:rsyslog]
command=service rsyslog start
priority=1
