# dockerfiles/dendrite/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -f /etc/ssh/sshd_config
directory=/
stderr_logfile=/var/log/ssh/sshd.log
stdout_logfile=/var/log/ssh/sshd_err.log
