# dockerfiles/dendrite/zookeeper.conf
[program:zookeeper]
command=service zookeeper-server start
priority=8
