# dockerfiles/dendrite/dendrite-elasticsearch/elasticsearch.conf
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch -Des.default.path.conf=/etc/elasticsearch
priority=1
