# dockerfiles/dendrite/dendrite-repos/tomcat.conf
[program:tomcat]
command=service tomcat7 start
priority=100
